# systemd/fix-aurdb.timer
[Unit]
Description=Rebuild local AUR DB on boot and periodically (system timer)

[Timer]
# Run 1 minute after boot to ensure DB exists
OnBootSec=1min
# Run hourly to catch any edge cases
OnUnitActiveSec=1h
Unit=fix-aurdb.service

[Install]
WantedBy=multi-user.target
