# systemd/fix-aurdb.path
[Unit]
Description=Watch local AUR repo directory and trigger DB rebuild

[Path]
PathExistsGlob=/var/lib/repo/aur/*.pkg.*
PathChanged=/var/lib/repo/aur
Unit=fix-aurdb.service

[Install]
WantedBy=multi-user.target
